# systemd/bot-engine.service
[Unit]
Description=Bot Engine v2.1 — Real-Time Scoring Server
After=network.target nginx.service
Wants=nginx.service
Documentation=https://github.com/your-org/bot-engine

[Service]
Type=simple
User=bot-engine
Group=bot-engine
ExecStart=/usr/bin/python3 /opt/bot-engine/realtime_server.py \
    --socket /run/bot-engine/scoring.sock
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=2
WatchdogSec=30

# Environment
# NOTE: BOT_HMAC_SECRET is NOT hard-coded here any more (C2.1). Put it
# in /etc/bot-engine/env with mode 0600, along with any other secrets.
# Example /etc/bot-engine/env:
#     BOT_HMAC_SECRET=<64-hex-char secret from `python3 -c 'import secrets;print(secrets.token_hex(32))'`>
# The engine will refuse to start if the secret is unset, a known
# placeholder, or shorter than 32 chars.
Environment=BOT_DB_PATH=/var/lib/bot-engine/bot_scores.db
Environment=BOT_BLOCKLIST_OUTPUT=/etc/nginx/conf.d/dynamic_blocklist.conf
Environment=BOT_HMAC_SECRET_FILE=/var/lib/bot-engine/.hmac_secret
# C4: admin endpoints and log level come from /etc/bot-engine/env by
# default (operators can drop the overrides there instead of editing
# this unit). Uncomment and set inline if you prefer unit-level config.
# Environment=BOT_ADMIN_ALLOW_IPS=127.0.0.1,::1,10.0.0.0/8
# Environment=BOT_LOG_LEVEL=INFO
EnvironmentFile=-/etc/bot-engine/env

# C4: SIGTERM graceful shutdown. The scoring server registers signal
# handlers that drain the DB write queue (see db_worker.py) and close
# the listener before exiting. 15 s is generous — worst case a full
# queue of 10k items takes ~1 s to flush on typical disks.
KillSignal=SIGTERM
TimeoutStopSec=15

# Runtime directory (auto-created, auto-cleaned)
RuntimeDirectory=bot-engine
RuntimeDirectoryMode=0755

# State directory for DB and model
StateDirectory=bot-engine
StateDirectoryMode=0750

# Logging — journald by default; set StandardOutput=file: for file-based logging
# StandardOutput=append:/var/log/bot-engine/scoring.log
# StandardError=append:/var/log/bot-engine/scoring.log
LogsDirectory=bot-engine
LogsDirectoryMode=0750

# Security hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/bot-engine /var/log/bot-engine /run/bot-engine /etc/nginx/conf.d
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
SystemCallFilter=@system-service
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
